# supabrain-worker.service
[Unit]
Description=SupaBrain Background Worker
Documentation=https://github.com/Scarface86c/supabrain
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/supabrain/core
ExecStart=/usr/bin/python3 /home/ubuntu/supabrain/core/supabrain_worker.py
Restart=always
RestartSec=10

# Environment
EnvironmentFile=/home/ubuntu/supabrain/.env

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=supabrain-worker

[Install]
WantedBy=multi-user.target
